# docs/systemd/hipsterpizza-cleanup.timer
[Unit]
Description=HipsterPizza weekly cleanup

[Timer]
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
